# sysutils/xe-guest-utilities/PLIST
@comment $NetBSD$
sbin/xe-daemon
sbin/xe-update-guest-attrs
share/examples/rc.d/xenguest
