# mail/opendkim/PLIST.filter
@comment $NetBSD$
man/man5/opendkim.conf.5
man/man8/opendkim.8
sbin/opendkim
share/examples/rc.d/opendkim
