# x11/xorg-server/PLIST.FreeBSD
@comment $NetBSD: PLIST.FreeBSD,v 1.1.1.1 2004/07/29 04:52:40 xtraeme Exp $
${X11ROOT_PREFIX}/etc/rc.d/xprint
${X11ROOT_PREFIX}/bin/kbd_mode
${X11ROOT_PREFIX}/lib/modules/freebsd/libdrm.a
${X11ROOT_PREFIX}/lib/modules/extensions/libdri.a
${X11ROOT_PREFIX}/${MAN_DIR}/kbd_mode.${MAN_SUFFIX}
@dirrm ${X11ROOT_PREFIX}/lib/modules/freebsd
@dirrm ${X11ROOT_PREFIX}/etc/rc.d
