# net/ra-rtsp-proxy/pkg/PLIST
@comment $NetBSD$
etc/rc.d/rtspd.sh
man/man8/rtspd.8
sbin/rtspd
share/doc/rtspd/README
share/doc/rtspd/license.txt
@dirrm share/doc/rtspd
