# net/iplog/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2003/04/06 11:06:04 hubertf Exp $
etc/rc.d/iplog
man/man5/iplog.conf.5
man/man8/iplog.8
sbin/iplog
share/examples/iplog/iplog.conf
@dirrm share/examples/iplog
