# mail/courier-mta/PLIST.ldap
@comment $NetBSD$
sbin/courierldapaliasd
share/examples/courier/ldapaliasrc.dist
share/examples/rc.d/courierldapaliasd
