# sysutils/pluged/PLIST
@comment $NetBSD$
sbin/pluged
share/examples/pluged/pluged.conf
share/examples/rc.d/pluged
