# security/cyrus-saslauthd/PLIST
@comment $NetBSD$
etc/rc.d/saslauthd
man/man8/saslauthd.8
sbin/saslauthd
share/examples/cyrus-saslauthd/LDAP_SASLAUTHD
@dirrm share/examples/cyrus-saslauthd
