# databases/postgresql-server/PLIST
@comment $NetBSD: PLIST,v 1.2 2002/02/15 10:12:34 skrll Exp $
bin/postgres
bin/postmaster
etc/rc.d/pgsql
man/man1/postgres.1
man/man1/postmaster.1
share/postgresql/postgres.bki
share/postgresql/postgres.description
share/postgresql/pg_hba.conf.sample
share/postgresql/pg_ident.conf.sample
share/postgresql/postgresql.conf.sample
@unexec ${RMDIR} -p %D/share/postgresql 2>/dev/null || ${TRUE}
