# chat/mu-conference/PLIST
@comment $NetBSD: PLIST,v 1.2 2005/05/02 20:33:58 reed Exp $
bin/mu-conference
share/examples/jabberd/mu-conference.sql
share/examples/jabberd/muc.xml
share/examples/jabberd/style.css
share/examples/rc.d/muc
