# devel/p4d/PLIST
@comment $NetBSD$
etc/rc.d/p4d
sbin/p4d
