# multimedia/minidlna/PLIST
@comment $NetBSD$
sbin/minidlna
share/examples/minidlna/minidlna.conf
share/examples/rc.d/minidlna
