# net/autonet/PLIST
@comment $NetBSD: PLIST,v 1.2 2003/05/29 11:45:50 abs Exp $
@dirrm share/examples/autonet
man/man8/autonet.8
sbin/autonet
share/examples/autonet/autonet.local
share/examples/autonet/lid_switch
share/examples/rc.d/autonet
