# mail/milter-spamc/PLIST
@comment $NetBSD$
sbin/milter-spamc
etc/rc.d/milter-spamc.sh
