# net/ipa/PLIST
@comment $NetBSD: PLIST,v 1.2 2002/02/15 10:12:54 skrll Exp $
bin/ipa
bin/ipastat
etc/rc.d/ipa
man/man5/ipa.5
man/man5/ipa.conf.5
man/man8/ipa.8
man/man8/ipastat.8
share/examples/ipa/README
share/examples/ipa/example-1
share/examples/ipa/example-2
share/examples/ipa/ipa.conf.default
@dirrm share/examples/ipa
