# sysutils/storage-manager/files/dptelog
#!/bin/sh
#
# $NetBSD$
#

# PROVIDE: dptelog
# REQUIRE: DAEMON dpteng

. /etc/rc.subr

name="dptelog"
rcvar=$name
command="@PREFIX@/lib/dpt/${name}"
start_cmd="dptelog_start"

dptelog_start()
{

	echo Starting $name.
	$command 2> /dev/null &
}

dptelog_stop()
{

	kill `ps ax | grep "lib/dpt/dptelog" | awk '{ print $1 }'`
}

load_rc_config $name
run_rc_command "$1" 
