# chat/i2cbd/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2002/11/14 00:59:31 hubertf Exp $
etc/rc.d/icbd
man/man8/icbd.8
man/man8/icbdb.8
sbin/icbd
sbin/icbdb
share/examples/i2cbd/empty.log
share/examples/i2cbd/icbd.deny
share/examples/i2cbd/icbd_help
share/examples/i2cbd/motd
share/examples/i2cbd/news.1
@dirrm share/examples/i2cbd
