# mail/milter-sender/PLIST
@comment $NetBSD$
sbin/milter-sender
etc/rc.d/milter-sender.sh
