# security/xdm-krb4/pkg/PLIST
@comment $NetBSD$
bin/xdm-krb4
etc/rc.d/xdm-krb4.sh
