# multimedia/ggrab/PLIST
@comment $NetBSD$
bin/ggrab
bin/sserver
share/doc/ggrab/README
share/examples/rc.d/sserver
@dirrm share/doc/ggrab
