# net/rbldnsd/PLIST
@comment $NetBSD$
${PKGMANDIR}/man8/rbldnsd.8
sbin/rbldnsd
share/doc/rbldnsd/NEWS
share/doc/rbldnsd/README.user
share/examples/rc.d/rbldnsd
@dirrm share/doc/rbldnsd
