# net/bpalogin/PLIST
@comment $NetBSD$
etc/rc.d/bpalogin
sbin/bpalogin
share/examples/bpalogin/bpalogin.conf
@dirrm share/examples/bpalogin
