# net/dhid/pkg/PLIST
@comment $NetBSD: PLIST,v 1.2 1999/08/04 14:04:37 jlam Exp $
sbin/dhid
sbin/dhid-genkeys
etc/rc.d/dhid
share/doc/dhid/COPYRIGHT
share/doc/dhid/README
share/doc/dhid/AUTHORS
share/doc/dhid/WHATSNEW
share/doc/dhid/INSTALL
@dirrm share/doc/dhid
