# databases/pgbouncer/PLIST
@comment $NetBSD$
bin/pgbouncer
man/man1/pgbouncer.1
man/man5/pgbouncer.5
share/examples/pgbouncer/pgbouncer.ini
share/examples/pgbouncer/pgbouncer.userlist
share/examples/rc.d/pgbouncer
