# net/delegate/pkg/PLIST
@comment $NetBSD$
sbin/delegated
etc/rc.d/delegated.sh.sample
man/cat8/delegated.0
share/doc/delegate/tutorial-en.html
share/doc/delegate/tutorial-jp.html
@dirrm share/doc/delegate
